# etc/profile.d/stuff.sh
export TERMINAL="alacritty"
export EDITOR="vim"
export OPENER="rifle"
eval "$(dircolors -b ~/.dircolors)"
